# addon/rootfs/etc/s6-overlay/s6-rc.d/init-packages/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Claude Code Add-on: Install additional packages
# ==============================================================================

if bashio::config.has_value 'additional_packages'; then
    bashio::log.info 'Installing additional packages...'

    apk update \
        || bashio::exit.nok 'Failed updating Alpine package indexes'

    for package in $(bashio::config 'additional_packages'); do
        apk add --no-cache "$package" \
            || bashio::exit.nok "Failed installing package: ${package}"
    done

    bashio::log.info 'Additional packages installed.'
fi
