# addon/rootfs/etc/s6-overlay/s6-rc.d/init-mcp/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Gemini CLI Add-on: Optional embedding model download
# ==============================================================================
declare enable_embeddings

enable_embeddings=$(bashio::config 'enable_embeddings')

if bashio::config.true 'enable_embeddings'; then
    bashio::log.info 'Embeddings enabled — downloading model (this may take a moment)...'

    ENABLE_EMBEDDINGS=true \
    MODELS_DIR=/data/models \
    node /opt/mcp-server/dist/warmup-model.js \
        || bashio::log.warning 'Embedding model download failed. Keyword search will still work.'

    bashio::log.info 'Embedding model ready.'
else
    bashio::log.info 'Embeddings disabled — using keyword-only search.'
fi
