# addon/rootfs/etc/s6-overlay/s6-rc.d/init-claude/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Claude Code Add-on: Environment setup
# Sets up API key, CLAUDE.md, .mcp.json, persistent storage
# ==============================================================================
declare api_key
declare model
declare yolo_mode
declare enable_embeddings

bashio::log.info 'Setting up Claude Code environment...'

# --- Read config ---
api_key=$(bashio::config 'api_key')
model=$(bashio::config 'model')
yolo_mode=$(bashio::config 'yolo_mode')
enable_embeddings=$(bashio::config 'enable_embeddings')

# --- Set environment variables for all services (s6 convention: one file per var) ---
printf '%s' "${enable_embeddings}" > /var/run/s6/container_environment/ENABLE_EMBEDDINGS
printf '%s' "/opt/mcp-server/data/docs-search.db" > /var/run/s6/container_environment/DOCS_DB_PATH
printf '%s' "/opt/docs" > /var/run/s6/container_environment/DOCS_PATH
printf '%s' "/data/models" > /var/run/s6/container_environment/MODELS_DIR

# Write all config to /etc/profile.d/ so the entrypoint shell sees them
{
    if bashio::var.has_value "${api_key}"; then
        echo "export ANTHROPIC_API_KEY=\"${api_key}\""
    else
        bashio::log.info 'No API key configured — use /login in Claude Code to authenticate via subscription.'
    fi
    echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\""
    echo "export DISABLE_AUTOUPDATER=1"
    echo "export CLAUDE_MODEL=\"${model}\""
    echo "export CLAUDE_YOLO=\"${yolo_mode}\""
} > /etc/profile.d/claude.sh

# --- Persistent storage ---
mkdir -p /data/.claude /data/models
if [[ ! -f /data/.claude.json ]]; then
    echo '{}' > /data/.claude.json
fi
touch /data/.bash_history

# Symlink persistent data into root home
ln -sf /data/.claude /root/.claude
ln -sf /data/.claude.json /root/.claude.json
ln -sf /data/.bash_history /root/.bash_history

# Symlink convenient paths
ln -sf /homeassistant /root/homeassistant
ln -sf /share /root/share

# --- Generate managed settings (system-wide) ---
mkdir -p /etc/claude-code

# When yolo mode is enabled, use managed settings to allow all tools.
# This works as root (unlike --dangerously-skip-permissions which blocks on root).
if bashio::config.true 'yolo_mode'; then
    bashio::log.info 'Yolo mode enabled — configuring auto-approve for all tools.'
    cat > /etc/claude-code/managed-settings.json << 'SETTINGSEOF'
{
  "permissions": {
    "allow": [
      "Bash",
      "Edit",
      "Write",
      "Read",
      "Glob",
      "Grep",
      "WebFetch",
      "WebSearch",
      "NotebookEdit",
      "TodoWrite",
      "Task",
      "EnterPlanMode",
      "ExitPlanMode",
      "SendMessage",
      "TeamCreate",
      "TeamDelete",
      "AskUserQuestion"
    ],
    "deny": []
  },
  "permissions.defaultMode": "acceptEdits"
}
SETTINGSEOF
else
    echo '{}' > /etc/claude-code/managed-settings.json
fi

# --- Generate managed MCP config (system-wide, always loaded) ---
bashio::log.info 'Generating MCP server configuration...'

cat > /etc/claude-code/managed-mcp.json << MCPEOF
{
  "mcpServers": {
    "home-assistant": {
      "command": "node",
      "args": ["/opt/mcp-server/dist/index.js"],
      "env": {
        "SUPERVISOR_TOKEN": "${SUPERVISOR_TOKEN}",
        "ENABLE_EMBEDDINGS": "${enable_embeddings}",
        "DOCS_DB_PATH": "/opt/mcp-server/data/docs-search.db",
        "DOCS_PATH": "/opt/docs",
        "MODELS_DIR": "/data/models"
      }
    }
  }
}
MCPEOF

# --- Generate default CLAUDE.md in /homeassistant/ if absent ---
if [[ ! -f /homeassistant/CLAUDE.md ]]; then
    bashio::log.info 'Creating default CLAUDE.md in Home Assistant config...'
    cat > /homeassistant/CLAUDE.md << 'EOF'
# Home Assistant Development Environment

You are working inside a Home Assistant instance with full API access.

## Available MCP Tools

The `home-assistant` MCP server provides these tools:
- `search_entities` - Search entities by name, domain, or area
- `get_entity_state` - Get state and attributes of an entity
- `call_service` - Call any HA service (turn on lights, trigger automations, etc.)
- `search_automations` - Search automation entities
- `get_ha_config` - Get HA core configuration
- `list_areas` - List all defined areas
- `search_devices` - Search device registry
- `get_config_entries` - List integration config entries
- `search_docs` - Search HA developer and user documentation
- `read_doc` - Read a specific documentation file
- `get_doc_stats` - Get documentation index statistics

## File Locations

- `/homeassistant/` - Home Assistant configuration directory
- `/homeassistant/.storage/` - Internal registry and state data
- `/share/` - Shared storage between add-ons
- `/ssl/` - SSL certificates (read-only)
- `/media/` - Media files (read-only)
EOF
fi

bashio::log.info 'Claude Code environment ready.'
